# ops/systemd/nullspace-simulator.service
[Unit]
Description=Nullspace Simulator API
After=network.target

[Service]
Type=simple
User=nullspace
Group=nullspace
WorkingDirectory=/opt/nullspace
EnvironmentFile=/etc/nullspace/simulator.env
ExecStart=/usr/local/bin/nullspace-simulator --host ${SIMULATOR_HOST} --port ${SIMULATOR_PORT} --identity ${SIMULATOR_IDENTITY} ${SIMULATOR_ARGS}
Restart=on-failure
RestartSec=5
LimitNOFILE=100000

[Install]
WantedBy=multi-user.target
