# ops/systemd/nullspace-consensus-watchdog.service
[Unit]
Description=Nullspace consensus watchdog
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/nullspace-consensus-watchdog.sh
User=root
Group=root

