# ops/systemd/nullspace-gateway.service
[Unit]
Description=Nullspace Gateway (WS)
After=network.target

[Service]
Type=simple
User=nullspace
Group=nullspace
WorkingDirectory=/opt/nullspace/gateway
EnvironmentFile=/etc/nullspace/gateway.env
Environment=NODE_ENV=production
ExecStart=/usr/bin/pnpm start
Restart=on-failure
RestartSec=5
LimitNOFILE=100000

[Install]
WantedBy=multi-user.target
